# PHP-HW1/myscript.service
[Unit]
Description=Мой PHP-скрипт
After=network.target

[Service]
Type=simple
ExecStart=/usr/bin/php /mnt/c/WORK/practicecode/php-cli/lesson-1/script.php >> /mnt/c/WORK/log2.txt
Restart=always
RestartSec=60

[Install]
WantedBy=multi-user.target